# docker/php/cron/schedule.cron
* * * * * /usr/bin/php8.0 /var/www/html/artisan schedule:run
